# roles/setup_check/files/systemd/checker@.service
# SPDX-FileCopyrightText: 2016-2026 Mirian Brechtel <markus.katharina.brechtel@thengo.net>
#
# SPDX-License-Identifier: EUPL-1.2

# checker systemd units

[Unit]
Description=Checker checks out %i

[Service]
WorkingDirectory=/etc/checker/checks/%i
ExecStart=/etc/checker/checker.sh /run/checker/%i.out
ExecStopPost=/etc/checker/notify.sh "%H" "%i" /run/checker/%i.out
Environment=PATH=/usr/lib/nagios/plugins:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/etc/checker/checks/%i/check.env
Type=exec
SuccessExitStatus=0 1 2 3

# Performance and resource limits (configurable via check.env)
TimeoutStartSec=${CHECK_TIMEOUT}
MemoryMax=${CHECK_MEMORY_MAX}
CPUQuota=${CHECK_CPU_QUOTA}

# Process limits (configurable via check.env)
LimitNOFILE=${CHECK_LIMIT_NOFILE}
LimitNPROC=${CHECK_LIMIT_NPROC}
