# roles/deploy_deploy/files/systemd/system/deploy@.path
# SPDX-FileCopyrightText: 2016-2026 Markus Katharina Brechtel <markus.katharina.brechtel@thengo.net>
#
# SPDX-License-Identifier: AGPL-3.0-or-later

[Unit]
Description=Deploy %i when deployment files change!

[Path]
PathChanged=/etc/deploy/%i/
TriggerLimitIntervalSec=60
TriggerLimitBurst=60
Unit=deploy@%i.service
