# roles/setup_deploy/files/systemd/system/deploy@.path
# SPDX-FileCopyrightText: 2016-2026 Mirian Brechtel <markus.katharina.brechtel@thengo.net>
#
# SPDX-License-Identifier: EUPL-1.2

[Unit]
Description=Deploy %i when deployment files change!

[Path]
PathChanged=/etc/deploy/%i/
TriggerLimitIntervalSec=60
TriggerLimitBurst=60
Unit=deploy@%i.service
