# test/quadlets/cute-devops-test@.container
# SPDX-FileCopyrightText: 2016 - 2026 Mirian Brechtel <markus.katharina.brechtel@thengo.net>
# SPDX-FileCopyrightText: 2020 - 2025 Uniklinik Köln
# SPDX-FileCopyrightText: 2025 - 2026 Goethe-University Frankfurt – Institute for Digital Medicine and Clinical Data Science
#
# SPDX-License-Identifier: EUPL-1.2

# Templated instance unit: one systemd system container per instance
# from the single cute-devops-test:trixie image. Start instances with
#   systemctl --user start cute-devops-test@mx
# The %i instance name becomes both the container name and its hostname,
# so peers resolve each other over the shared network by that name.
#
# This is the working prototype of the deploy_quadlet role proposed in
# issues/container-apps.feature.md — a rootless, linger-backed, quadlet
# system container.
[Unit]
Description=cute-devops test system container (%i)

[Container]
Image=cute-devops-test:trixie
ContainerName=%i
HostName=%i
Network=cute-devops-test.network
# podman auto-detects the /sbin/init entrypoint and wires up /run, /tmp
# and /sys/fs/cgroup for systemd; be explicit so it holds regardless of
# how the unit is invoked.
# --systemd=always: run systemd as PID 1 (tmpfs /run, cgroup wiring).
# --no-hosts: don't bind-mount /etc/hosts, so the hostname role can write
# it in place; cross-instance name resolution still works via aardvark DNS.
PodmanArgs=--systemd=always --no-hosts
# No host PublishPort by default: rootless cannot bind <1024 and the
# end-to-end test runs inside the shared network. Uncomment to poke an
# instance from the host on a high port, e.g. mb's IMAPS:
#PublishPort=1993:993

[Service]
# A test box should not silently resurrect mid-run; restart is manual.
Restart=no
TimeoutStartSec=180

[Install]
WantedBy=default.target
